# project-aria-bot.service
[Unit]
After=postgresql.service

[Service]
ExecStart=/usr/bin/init_bot.sh

[Install]
WantedBy=default.target

